# afanasy/package/afanasy-server/DEBIAN/preinst
#!/bin/bash

echo "Afanasy server PRE INSTALL"

systemctl stop afserver || true
systemctl disable afserver || true
systemctl daemon-reload || true

exit 0
