# afanasy/package/afanasy-server/DEBIAN/prerm
#!/bin/bash

echo "Afanasy server PRE REMOVE"

systemctl stop afserver || true
systemctl disable afserver || true
systemctl daemon-reload || true

exit 0
