# afanasy/package/afanasy-render/DEBIAN/postinst
#!/bin/bash

echo "Afanasy render POST INSTALL"

if ! id render; then
	if getent group render; then
		useradd render -g render -m
	else
		useradd render -m
	fi
fi

systemctl enable /opt/cgru/afanasy/systemd/afrender.service || true
systemctl daemon-reload || true
systemctl start afrender.service || true

exit 0

