# afanasy/package/afanasy-render/DEBIAN/preinst
#!/bin/bash

echo "Afanasy render PRE INSTALL"

systemctl stop afrender || true
systemctl disable afrender || true
systemctl daemon-reload || true

exit 0
